# leonardo/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:{{ service.type }}_{{ service.name }}]
command=/srv/leonardo/sites/{{ service.name }}/bin/gunicorn_start
stdout_logfile=/srv/leonardo/sites/{{ service.name }}/logs/access.log 
stderr_logfile=/srv/leonardo/sites/{{ service.name }}/error.log
user=leonardo
autostart=true
autorestart=true
