# leonardo/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- set app = salt['pillar.get']('leonardo:server:app:'+service.name) %}

[program:{{ service.type }}_{{ service.name }}]
command=/srv/leonardo/sites/{{ service.name }}/bin/gunicorn_start
stdout_logfile=/srv/leonardo/sites/{{ service.name }}/logs/access.log 
stderr_logfile=/srv/leonardo/sites/{{ service.name }}/error.log
user=leonardo
autostart=true
autorestart=true

{% if app.worker is defined and app.worker %}

[program:{{ service.type }}_{{ service.name }}_worker]
directory=/srv/leonardo/sites/{{ service.name }}
command=/srv/leonardo/sites/{{ service.name }}/bin/python /srv/leonardo/sites/{{ service.name }}/manage.py celery worker -B -E -Q {{ service.name }}
stdout_logfile=/srv/leonardo/sites/{{ service.name }}/logs/worker.log 
stderr_logfile=/srv/leonardo/sites/{{ service.name }}/worker_error.log
user=leonardo
autostart=true
autorestart=true

[program:{{ service.type }}_{{ service.name }}_celerycam]
directory=/srv/leonardo/sites/{{ service.name }}
command=/srv/leonardo/sites/{{ service.name }}/bin/python /srv/leonardo/sites/{{ service.name }}/manage.py celerycam
stdout_logfile=/srv/leonardo/sites/{{ service.name }}/logs/celerycam.log 
stderr_logfile=/srv/leonardo/sites/{{ service.name }}/celerycam_error.log
user=leonardo
autostart=true
autorestart=true

{% endif %}
